# main/postmarketos-bootsplash/kill-pbsplash.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name="Kill pbsplash"
description="Kill pbsplash"
command="killall pbsplash"

depend() {
	after agetty.tty1
}
